# flask-blog.service
[Unit]
Description=Flask Application (point2repo)
[Service]
User=pythonadm
WorkingDirectory=/opt/flask-blog
ExecStart=/bin/bash \
	-c 'cd /opt/flask-blog && \
	source ./venv/bin/activate && \
	python ./app.py'
Environment=SECRET_KEY={{SECRET_KEY}}
Environment=DATABASE_URI=mysql+pymysql://{{MYSQL_USER}}:{{MYSQL_PASSWORD}}@{{MYSQL_HOST}}/{{MYSQL_DATABASE}}
[Install]
WantedBy=multi-user.target
